# deploy/soul-v2-tasks.service
[Unit]
Description=Soul v2 — Tasks Server
After=network.target

[Service]
Type=simple
User=rishav
Group=rishav
WorkingDirectory=/home/rishav/soul-v2
ExecStart=/home/rishav/soul-v2/soul-tasks serve
Restart=on-failure
RestartSec=5
Environment=SOUL_TASKS_HOST=127.0.0.1
Environment=SOUL_TASKS_PORT=3004

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/rishav/.soul-v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
