# deploy/soul-v2-scout.service
[Unit]
Description=Soul v2 — Scout Server
After=network.target

[Service]
Type=simple
User=rishav
Group=rishav
WorkingDirectory=/home/rishav/soul-v2
ExecStart=/home/rishav/soul-v2/soul-scout serve
Restart=on-failure
RestartSec=5
Environment=SOUL_SCOUT_HOST=127.0.0.1
Environment=SOUL_SCOUT_PORT=3020

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/rishav/.soul-v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
