# deploy/soul-v2-infra.service
[Unit]
Description=Soul v2 — Infra Server
After=network.target

[Service]
Type=simple
User=rishav
Group=rishav
WorkingDirectory=/home/rishav/soul-v2
ExecStart=/home/rishav/soul-v2/soul-infra serve
Restart=on-failure
RestartSec=5
Environment=SOUL_INFRA_HOST=127.0.0.1
Environment=SOUL_INFRA_PORT=3012

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/rishav/.soul-v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
