# deploy/soul-v2-docs.service
[Unit]
Description=Soul v2 — Docs Server
After=network.target

[Service]
Type=simple
User=rishav
Group=rishav
WorkingDirectory=/home/rishav/soul-v2
ExecStart=/home/rishav/soul-v2/soul-docs serve
Restart=on-failure
RestartSec=5
Environment=SOUL_DOCS_HOST=127.0.0.1
Environment=SOUL_DOCS_PORT=3018

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/rishav/.soul-v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
