# deploy/soul-v2-bench.service
[Unit]
Description=Soul v2 — Bench Server
After=network.target

[Service]
Type=simple
User=rishav
Group=rishav
WorkingDirectory=/home/rishav/soul-v2
ExecStart=/home/rishav/soul-v2/soul-bench serve
Restart=on-failure
RestartSec=5
Environment=SOUL_BENCH_HOST=127.0.0.1
Environment=SOUL_BENCH_PORT=3026

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/rishav/.soul-v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
